# etc/nix/home/burke.gitpull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key>
 <string>burke.gitpull</string>
 <key>ProgramArguments</key>
 <array>
   <string>/Users/burke/.nix-profile/bin/launchd-gitpull</string>
 </array>
 <key>StandardErrorPath</key>
 <string>/tmp/burke.gitpull.err</string>
 <key>StandardOutPath</key>
 <string>/tmp/burke.gitpull.out</string>
 <key>StartInterval</key>
 <integer>20</integer>
</dict>
</plist>
